# vpn
#!/bin/sh /etc/rc.common
START=99
STOP=15
 
start() {
  ulimit -SHn 65535
  ss restart
}
 
stop() {
  ss stop
}
