# vpn
#!/bin/sh /etc/rc.common
START=99
STOP=15
 
start() {
  ss start
}
 
stop() {
  ss stop
}
